# n64.service
[Unit]
Description=N64 Controller
Before=systemd-sysctl.service

[Service]
Type=oneshot
RemainAfterExit=yes
KillMode=none
ExecStart=/opt/n64/start.sh

[Install]
WantedBy=multi-user.target
